# bspwm/global/etc/profile.d/bspwm_panel.sh
#!/bin/sh
export PATH="$PATH:/usr/local/bin/panel"

export PANEL_FIFO=/tmp/panel-fifo
export PANEL_HEIGHT=24
export PANEL_FONT_FAMILY="-*-terminus-medium-r-normal-*-12-*-*-*-c-*-*-1"
