# bspwm/global/etc/profile.d/bspwm_panel.sh.MajesticMoose.local
#!/bin/sh
export PANEL_FONT_FAMILY="LiberationMono:size=09"
