# systemd/home/.config/systemd/user/rotatebackground.service
[Unit]
Description=Rotate system background

[Service]
Environment=DISPLAY=:0
ExecStart=/home/jamie/bin/rotatebackground.sh
ExecStop=/bin/true
