# bspwm/global/etc/system/bluetooth.service.d/override.conf
[Service]
ExecStart=/usr/lib/bluetooth/bluetoothd -E
